# packages/bitcore-node/bitcore-node.service
[Unit]
Description=Bitcore Node
After=network.target

[Service]
User=ubuntu
Group=ubuntu
ExecStart=/usr/bin/node /home/ubuntu/coins/Mtnscore/packages/bitcore-node/build/src/server.js
WorkingDirectory=/home/ubuntu/coins/be

[Install]
WantedBy=multi-user.target
